# vagrant/share/build/docker/nginx.service
[Unit]
Description=nginx
Requires=docker.service
After=docker.service
[Service]
RemainAfterExit=yes
Restart=always
ExecStart=/usr/bin/docker start docker_nginx_1
ExecStop=/usr/bin/docker stop docker_nginx_1
PIDFile=/tmp/pids/docker_nginx_1.pid
TimeoutSec=30
[Install]
WantedBy=multi-user.target
 